# systemd/user/udiskie.service
[Unit]
Description=Udiskie automount

[Service]
ExecStart=/usr/bin/udiskie -aNC2

[Install]
WantedBy=default.target
